# docker/template/golang/adapters/conf/supervisor.d/adapters.ini
[program:adapters]
command=/tiup/adapters
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/tiup/logs/supervisor/adapters.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/tiup/logs/supervisor/adapters.err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
