# docker/template/tiup_file/edu_mp4/conf/supervisor.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm -c /tiup/conf/php-fpm/php.ini -y /tiup/conf/php-fpm/php-fpm.conf
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/tiup/logs/supervisor/php-fpm.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/tiup/logs/supervisor/php-fpm.err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
