# docker/template/ice/tagsearch/conf/supervisor.d/dlogd.ini
[program:dlogd]
command=/tiup/dlogd -l /tiup/dlogd.plugin -D
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/tiup/logs/supervisor/dlogd.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/tiup/logs/supervisor/dlogd.err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
