# docker/template/tiup_server/redis/conf/supervisor.d/redis.ini
[program:redis]
command=/usr/sbin/redis-server /tiup/conf/redis/redis.conf
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/tiup/logs/supervisor/redis.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/tiup/logs/supervisor/redis.err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
