# docker/template/tiup_file/pic-img/conf/supervisor.d/rsyncd.ini
[program:rsyncd]
command=/usr/bin/rsync --daemon --config=/tiup/conf/rsync/rsyncd.conf --no-detach
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/tiup/logs/supervisor/sshd.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/tiup/logs/supervisor/sshd.err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
